# role/redis/templates/supervisor.ini.j2
[program:redis]
command=bash -c "redis-server {{server_conf_path}}"
process_name=%(program_name)s
pidfile={{server_pid_path}}
autorestart=false
