# role/redis/templates/supervisor.ini.j2
[program:redis]
command=bash -c "redis-server {{redis.conf_dir}}/redis.conf"
process_name=%(program_name)s
pidfile={{redis.system_dir}}/process.pid
autorestart=false
