# role/nginx/templates/supervisor.ini.j2
[program:nginx]
command=bash -c "nginx"
process_name=%(program_name)s
pidfile={{server_pid_path}}
autorestart=false
