# supervisor/conf/sshd.ini
[program:sshd]
command=bash -c "/usr/sbin/sshd -D"
process_name=%(program_name)s
autorestart=false