# roles/redis/templates/redis.ini.j2
[program:redis]
command=bash -c "/usr/local/bin/redis-server {{confdir}}/redis.conf"
process_name=%(program_name)s
pidfile={{ systemdir }}/process.pid
autorestart=false