# playbooks/roles/redis/templates/supervisor.ini.j2
[program:redis]
command=bash -c "/usr/local/bin/redis-server {{redis.config_dir}}/redis.conf"
process_name=%(program_name)s
pidfile={{redis.system_dir}}/process.pid
autorestart=false
autostart=false
